# deploy/systemd/cubit-backend.service
# CUBIT auth backend (Bun + SQLite). Install:
#   sudo cp deploy/systemd/cubit-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now cubit-backend
# Adjust User/WorkingDirectory/paths to match your server.
[Unit]
Description=CUBIT auth backend (Bun + SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cubit/backend
# EnvironmentFile holds CUBIT_SECRET, CUBIT_GAME_URL (see deploy/.env.example)
EnvironmentFile=/opt/cubit/deploy/.env
# Port pinned here (NOT in .env, which the game server also loads).
Environment=PORT=3001
ExecStart=/root/.bun/bin/bun server.ts
Restart=always
RestartSec=3
# SQLite data lives in WorkingDirectory (data.sqlite). Keep it writable + backed up.

[Install]
WantedBy=multi-user.target
