# deploy/systemd/cubit-game.service
# CUBIT game server (HYTOPIA SDK). Install:
#   sudo cp deploy/systemd/cubit-game.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now cubit-game
#
# Builds the .mjs with bun (ExecStartPre) but RUNS it on Node: the WebTransport
# native QUIC module crashes Bun (Illegal instruction) the moment a client
# connects — Node is stable. Needs toktx on PATH + LD_LIBRARY_PATH to build the
# block-texture atlas at startup, else it crashes with `spawn toktx ENOENT`.
[Unit]
Description=CUBIT game server (HYTOPIA)
After=network-online.target cubit-backend.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cubit
EnvironmentFile=/opt/cubit/deploy/.env
# The game server also reads PORT; pin it to 8080 (the .env no longer sets PORT).
Environment=PORT=8080
# toktx (installed by setup.sh under $KTX_HOME=/opt/ktx)
Environment=PATH=/root/.bun/bin:/opt/ktx/bin:/usr/local/bin:/usr/bin:/bin
Environment=LD_LIBRARY_PATH=/opt/ktx/lib
# Build with bun, run with node (Bun crashes on WebTransport connect — see above).
ExecStartPre=/root/.bun/bin/bun run node_modules/.bin/hytopia build index.ts
ExecStart=/usr/bin/node --enable-source-maps /opt/cubit/index.mjs
Restart=always
RestartSec=5
# The SDK serves HTTPS + WebTransport directly on :8080 — open TCP+UDP 8080.

[Install]
WantedBy=multi-user.target
